# docker-snapshot/crontab

# Every 10mins from 7am to 7pm MST
*/10 13-23,00,01 * * * /data/make_image.sh
